# lib/svc/manifest/snltd_monitor.xml
<?xml version='1.0'?>

<!DOCTYPE service_bundle SYSTEM '/usr/share/lib/xml/dtd/service_bundle.dtd.1'>

<service_bundle type='manifest' name='export'>

	<service name='snltd/monitor' type='service' version='1'>

		<create_default_instance enabled='true'/>

		<single_instance/>

		<dependency 
			name='milestone' 
			grouping='require_all'
			restart_on='error'
			type='service'>
			<service_fmri value='svc:/milestone/multi-user:default'/>
		</dependency>

		<exec_method
			name='start'
			type='method'
			exec='/usr/local/snltd_monitor/bin/snltd_monitor_daemon.sh -D'
			timeout_seconds='60'>

			<method_context>
				<method_credential
					user='monitor'
					group='other'
					privileges='basic,file_dac_read,file_dac_search,!file_link_any'/>
			</method_context>

		</exec_method>

		<exec_method
			name='stop'
			type='method'
			exec=':kill'
			timeout_seconds='60'>
			</exec_method>

		<stability value='Evolving'/>

		<template>
			<common_name>
				<loctext xml:lang='C'>SNLTD host monitor daemon</loctext>
			</common_name>

			<documentation>
				<doc_link name='GitHub' uri='https://github.com/snltd/snltd_monitor' /> 
			</documentation>

		</template>

	</service>

</service_bundle>

